# tasks.cron
##
# This file is automatically imported into cronjobs on container startup.
# IMPORTANT NOTE!
# Please make sure there is a blank line after the last cronjob entry.
##

# Use this if you want to use real cron instead of wp-cron
# You need to have define('DISABLE_WP_CRON', true); in  your wp-config.php
# This is the default way of bedrock
* * * * * /var/www/project/scripts/run-cron.sh
